# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Admin dashboard — always on. Its Start/Stop buttons drive the gateway via
# supervisorctl (GATEWAY_SUPERVISOR=gateway).
[program:dashboard]
command=/usr/bin/python3 /home/it-admin/wa-lead-gen/admin/app.py
directory=/home/it-admin/wa-lead-gen
environment=HOME="/home/it-admin",GATEWAY_SUPERVISOR="gateway",SUPERVISOR_CONF="/home/it-admin/wa-lead-gen/docker/supervisord.conf",ADMIN_HOST="0.0.0.0",ADMIN_PORT="8088",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# openclaw gateway — auto-starts with the container; the dashboard can stop/start it.
# Logs to progress/gateway.log (the dashboard reads this for status).
[program:gateway]
command=/usr/local/bin/openclaw gateway
directory=/home/it-admin/wa-lead-gen
environment=HOME="/home/it-admin",PATH="/usr/local/bin:/usr/bin:/bin",OPENCLAW_AUTO_UPDATE="0"
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/home/it-admin/wa-lead-gen/progress/gateway.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/home/it-admin/wa-lead-gen/progress/gateway.log
priority=20
